# cm-1/zabbix-checks.cron
*/4 * * * * root tail -n 120 /var/log/ceph/ceph.log | grep degraded; zabbix_sender -z 10.1.2.217 -s $(hostname) -k ceph-degraded -o $? > /dev/null
