# shared/clean-cache.cron
MAILTO=""

# Clean the apt cache immediately before the weekly backup job to reduce storage overhead
45 3 * * 1 root bash -c "apt-get clean all"

